# deploy/systemd/invoice-backend.service
[Unit]
Description=Invoice FastAPI Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/invoice/backend
EnvironmentFile=/opt/invoice/backend/.env
ExecStart=/opt/invoice/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
